# scripts/com.leaderboard.swim-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Daily swim-results poll. Install with:
       cp scripts/com.leaderboard.swim-poll.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.leaderboard.swim-poll.plist
     Runs from the main checkout. Edit ProgramArguments path only if the repo moves.
     Logs: ~/Library/Logs/leaderboard-swim-poll.log -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leaderboard.swim-poll</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/clawdnett/Projects/leaderboard/scripts/swim-poll.sh</string>
    </array>

    <!-- Run once a day at 07:00 local. launchd also fires a missed run when the
         Mac next wakes, so a sleeping mini still catches up. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/clawdnett/Library/Logs/leaderboard-swim-poll.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/clawdnett/Library/Logs/leaderboard-swim-poll.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
